# pi-para-daemon.service
[Unit]
Description=pi-para knowledge capture daemon
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/projects/pi-para
ExecStart=/home/%i/.local/share/mise/installs/node/22.22.0/bin/npx tsx src/cli.ts start
Restart=on-failure
RestartSec=10
Environment=HOME=/home/%i
Environment=PATH=/home/%i/.local/share/mise/installs/node/22.22.0/bin:/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
