# pi-para-daemon.service
[Unit]
Description=pi-para knowledge capture daemon
After=network.target

[Service]
Type=simple
User=%i
# Use login shell to inherit the user's PATH (picks up mise, nvm, fnm, volta, or system node).
# -l sources ~/.profile / ~/.bash_profile which set up node version managers.
ExecStart=/bin/bash -lc 'exec npx pi-para-daemon start'
Restart=on-failure
RestartSec=10
Environment=HOME=/home/%i

[Install]
WantedBy=multi-user.target
